# ansible/roles/deploy/templates/chemvr.service.j2
[Unit]
Description=ChemVR Application
After=network.target

[Service]
Type=simple
User={{ app_user }}
WorkingDirectory={{ app_dir }}
Environment="NODE_ENV=production"
Environment="PORT={{ app_port }}"
# Add other environment variables here
ExecStart=/usr/bin/pnpm start
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=chemvr

[Install]
WantedBy=multi-user.target
